# development/crontab
* * * * * sched scheduler c "echo Have fun!"